# sv/barbican/run
#!/bin/sh

#sv start socklog-unix || exit 1

export VIRTUAL_ENV=/myenv

exec 2>&1
exec barbican.sh start
